# debian/postinst
#!/bin/sh
# postinst script for ftp-cloudfs
#
# see: dh_installdeb(1)

set -e

# summary of how this script can be called:
#        * <postinst> `configure' <most-recently-configured-version>
#        * <old-postinst> `abort-upgrade' <new version>
#        * <conflictor's-postinst> `abort-remove' `in-favour' <package>
#          <new-version>
#        * <postinst> `abort-remove'
#        * <deconfigured's-postinst> `abort-deconfigure' `in-favour'
#          <failed-install-package> <version> `removing'
#          <conflicting-package> <version>
# for details, see http://www.debian.org/doc/debian-policy/ or
# the debian-policy package


case "$1" in
    configure)
	if ! getent group ftpcloudfs > /dev/null 2>&1
	then
		addgroup --system ftpcloudfs >/dev/null
	fi
	if ! getent passwd ftpcloudfs > /dev/null 2>&1
	then
	        adduser --system --home /var/log/ftp-cloudfs --ingroup ftpcloudfs --no-create-home --shell /bin/false ftpcloudfs > /dev/null 2>&1
	fi

        # log directory
	mkdir -m 0700 -p /var/log/ftp-cloudfs
	touch /var/log/ftp-cloudfs/ftp-cloudfs.log
	chmod 0640 /var/log/ftp-cloudfs/ftp-cloudfs.log
	chown -R ftpcloudfs:ftpcloudfs /var/log/ftp-cloudfs

	# pid file directory
	mkdir -m 0700 -p /var/run/ftp-cloudfs
	chown -R ftpcloudfs:ftpcloudfs /var/run/ftp-cloudfs

	if which update-python-modules >/dev/null 2>&1; then
        	update-python-modules --post-install
	fi

	if [ -x "/etc/init.d/ftp-cloudfs" ]; then
        	update-rc.d ftp-cloudfs defaults >/dev/null
        #	if [ -x "`which invoke-rc.d 2>/dev/null`" ]; then
        #        	invoke-rc.d ftp-cloudfs start || exit $?
        #	else
        #        	/etc/init.d/ftp-cloudfs start || exit $?
        #	fi
	fi

	echo "*******************************************"
	echo "Your UID is:"
	id -u ftpcloudfs
	echo "Your GID is:"
	id -g ftpcloudfs
	echo "Update your configuration file accordingly."
	echo "*******************************************"

    ;;

    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
    ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0
